# services/badvpn-udpgw.service
[Unit]
Description=BadVPN UDPGW - UDP Gateway para VPS-MX
Documentation=https://github.com/nicholasgasior/badvpn
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/badvpn-udpgw --listen-addr 127.0.0.1:7300 --max-clients 1000 --max-connections-for-client 10
Restart=always
RestartSec=3
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
